# config/etc/systemd/system/suspend-server.service
[Unit]
Description="suspend server (after checks)"

[Service]
ExecStart=/home/marc/bin/suspend-server
